# SRV1C_host/srv1cv83.service
[Unit]
Description=1C:Enterprise server
After=network.target
[Service]
Type=simple
ExecStart=/opt/1C/v8.3/x86_64/ragent -d /_data/srv1c_inf_log -debug
KillSignal=SIGINT
PrivateTmp=true
Restart=on-failure
RestartSec=5
User=usr1cv8
Group=grp1cv8
[Install]
WantedBy=multi-user.target
